# deploy/accurate-giant.supervisor.conf
; Supervisor program for the Wobedi Bi Lotto Next.js production server.
;
; NOTE: The program name and on-server paths still read `accurate-giant`
; because they're tied to the live deployment. Rename in lockstep with a
; fresh cutover (new path layout under /home/admin/wobedibi/, new
; supervisor program name, new nginx vhost slug).
;
; Runs the standalone build via /home/admin/accurate-giant/start.sh, which
; loads the .env and exec's `node server.js`. Listens on 127.0.0.1:3002 and
; is reverse-proxied by /etc/nginx/sites-available/accurate-giant.

[program:accurate-giant]
process_name=%(program_name)s
directory=/home/admin/accurate-giant/current
command=/home/admin/accurate-giant/start.sh
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=admin
environment=NODE_ENV="production",PORT="3002",HOSTNAME="127.0.0.1"
stdout_logfile=/var/log/supervisor/accurate-giant.log
stderr_logfile=/var/log/supervisor/accurate-giant.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
